# .config/systemd/user/ut-moosicmaan-backup.service
[Unit]
Description=Moosicmaan Backup

[Service]
Type=oneshot
ExecStart=/home/moosicmaan/.config/.scripts/ut-moosicmaan-backup

