# etc/systemd/system/mnt-backup.automount
[Unit]
Description=Toshiba External USB Backup Drive

[Automount]
Where=/mnt/backup
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
