# etc/systemd/system/quote-broadcast.service
[Unit]
Description=Broadcast a random quote to all users
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/broadcast-quote.sh
